# logging/appinsightslogger.conf

description "application insights logger"

start on runlevel [2345]
stop on runlevel [!2345]

exec sudo -u serveruser python /opt/conf/logging/appinsightslogger.py >$

respawn

respawn limit 10 5
